# e.conf
[program:e]
command=python3 /home/curator2/e/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/e/test.err.log
stdout_logfile=/home/curator2/e/test.out.log